# deploy/systemd/camstation-updater.service
[Unit]
Description=CamStation auto-updater

[Service]
Type=oneshot
ExecStart=/opt/camstation/deploy/deploy.sh
StandardOutput=journal
StandardError=journal
